# keepalived-vip/init/keepalived.conf
start on stopped [12345]
stop on runlevel [!12345]

script
    . /etc/sysconfig/keepalived
    exec /sbin/keepalived $KEEPALIVED_OPTIONS

end script

respawn
